# web-api/officelink.ini
[program:officelink]
command=/usr/local/officelink/bin/python -u /var/www/python/officelink.py
autorstart=true
